# tyrion/tyrion-node.init
#!/bin/sh
#
# tyrion-node  init file for starting up the Tyrion node
#
# chkconfig:   - 20 80
# description: Starts and stops the Tyrion daemon that handles all \
#              service requests.

# Source function library.
. /etc/rc.d/init.d/functions

exec="/usr/sbin/tyrion-node-wrapper"
config="/etc/tyrion/node.conf"
logfile="/var/log/tyrion/node.out.log"

[ -e /etc/sysconfig/tyrion-node ] && . /etc/sysconfig/tyrion-node

lockfile=/var/lock/subsys/tyrion-node

start() {
    [ -x $exec ] || exit 5
    echo -n $"Starting tyrion-node: "
    daemon --user root "$exec -c $config >> $logfile 2>&1 &"
    retval=$?
    echo
    [ $retval -eq 0 ] && touch $lockfile
    return $retval
}

stop() {
    echo -n $"Stopping tyrion-node: "
    killproc tyrion-node-wrapper
    retval=$?
    echo
    [ $retval -eq 0 ] && rm -f $lockfile
    return $retval
}

restart() {
    stop
    start
}

reload() {
    killproc tyrion-node -HUP
}

rh_status() {
    status tyrion-node-wrapper
}

rh_status_q() {
    rh_status >/dev/null 2>&1
}


case "$1" in
    start)
        rh_status_q && exit 0
        $1
        ;;
    stop)
        rh_status_q || exit 0
        $1
        ;;
    restart)
        $1
        ;;
    reload)
        rh_status_q || exit 7
        $1
        ;;
    force-reload)
        force_reload
        ;;
    status)
        rh_status
        ;;
    condrestart|try-restart)
        rh_status_q || exit 0
        restart
        ;;
    *)
        echo $"Usage: $0 {start|stop|status|restart|condrestart|try-restart|reload}"
        exit 2
esac
exit $?
